# etc-systemd-system-myapp.service/myapp.service
[Unit]
Description=My Flask App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/myapp
ExecStart=/usr/bin/python3 /root/myapp/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
